# docker/crontab
# Bird of the Day — supercronic crontab.
#
# Daily generation at 07:00 UTC, matching the existing GitHub Actions
# workflow (.github/workflows/ave-del-dia.yml). To change the schedule,
# edit this file and rebuild the image (or mount a custom one at
# /etc/supercronic/crontab).
#
# Format: standard crontab. supercronic logs each run to stdout, so the
# Python script's logging falls through to `docker logs`.

0 7 * * * /opt/venv/bin/python -m scripts.generate
